# services/snakeeye-web.service
[Unit]
Description=SnakeEye Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/snakeeye/www
ExecStart=/usr/bin/python3 -m http.server 8000 --bind 0.0.0.0
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
